# deploy/systemd/face-recognition.service
# systemd-юнит: распознавание лиц + ANPR (пишет в БД/галерею).
# Установка:
#   sudo cp deploy/systemd/face-recognition.service /etc/systemd/system/
#   отредактируй User и пути под свой сервер, затем:
#   sudo systemctl daemon-reload && sudo systemctl enable --now face-recognition
[Unit]
Description=Face + ANPR recognition service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# путь, куда склонирован/скопирован проект:
WorkingDirectory=/root/inspeksiya-face-recognition
Environment=PYTHONUNBUFFERED=1
Environment=INSIGHTFACE_HOME=/root/inspeksiya-face-recognition/.insightface
# gpu_setup сам предзагружает CUDA .so из venv (nvidia-*-cu12); LD_LIBRARY_PATH не обязателен.
ExecStart=/root/inspeksiya-face-recognition/.venv/bin/python src/main.py
Restart=always
RestartSec=5
# не убивать по OOM слишком агрессивно
OOMScoreAdjust=-100

[Install]
WantedBy=multi-user.target
