# deploy/systemd/face-dashboard.service
# systemd-юнит: FastAPI дашборд (читает БД/файлы, удаляет записи).
# Установка:
#   sudo cp deploy/systemd/face-dashboard.service /etc/systemd/system/
#   отредактируй User и пути, затем:
#   sudo systemctl daemon-reload && sudo systemctl enable --now face-dashboard
[Unit]
Description=Face + ANPR dashboard (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/inspeksiya-face-recognition
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/inspeksiya-face-recognition/.venv/bin/python3 -m uvicorn web.app:app --host 0.0.0.0 --port 8089 --app-dir /root/inspeksiya-face-recognition
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
